# init/ubuntu/azure-vmextensions.slice
[Unit]
Description=Slice for Azure VM Extensions
